# lined/lined.service
[Unit]
Description=get lined up

[Service]
Restart=always
ExecStart=/usr/local/bin/lined

[Install]
WantedBy=multi-user.target
